# crontab
*/15 * * * * /usr/local/bin/php /var/www/update.php --feeds
